# envoyou-task-processor.service
[Unit]
Description=Envoyou Background Task Processor
After=network.target
Requires=network.target

[Service]
Type=simple
User=envoyou
Group=envoyou
WorkingDirectory=/home/husni/PROJECT-ENVOYOU-API/api-envoyou
Environment=PATH=/home/husni/PROJECT-ENVOYOU-API/api-envoyou/venv/bin
Environment=PYTHONPATH=/home/husni/PROJECT-ENVOYOU-API/api-envoyou
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/husni/PROJECT-ENVOYOU-API/api-envoyou/venv/bin/python run_task_processor.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=envoyou-task-processor

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/husni/PROJECT-ENVOYOU-API/api-envoyou
ProtectHome=yes

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target